# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=5

[program:php-fpm]
command=php-fpm --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
priority=10
user=www-data